# etc/systemd/system/bluetooth.service.d/override.conf
# Disable the SAP plugin to remove error
[Service]
ExecStart=
ExecStart=/usr/libexec/bluetooth/bluetoothd --noplugin=sap